# deployment/services/launchd/com.contractanalyzer.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.contractanalyzer.frontend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /opt/contract-analyzer/frontend && source ../venv/bin/activate && exec streamlit run app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --server.enableCORS false --server.enableXsrfProtection true</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/opt/contract-analyzer/frontend</string>
    
    <key>UserName</key>
    <string>contract-analyzer</string>
    
    <key>GroupName</key>
    <string>contract-analyzer</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/contract-analyzer/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/opt/contract-analyzer</string>
    </dict>
    
    <key>EnvironmentVariablesFile</key>
    <string>/etc/contract-analyzer/environment</string>
    
    <key>StandardOutPath</key>
    <string>/var/log/contract-analyzer/frontend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/contract-analyzer/frontend-error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LowPriorityIO</key>
    <false/>
    
    <key>Nice</key>
    <integer>0</integer>
    
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
        <key>NumberOfProcesses</key>
        <integer>2048</integer>
    </dict>
    
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
        <key>NumberOfProcesses</key>
        <integer>2048</integer>
    </dict>
</dict>
</plist>